# roles/letsencrypt/templates/ssl-renew.service
[Unit]
Description=Renew SSL Certs

[Service]
ExecStart=/usr/bin/certbot -c /etc/letsencrypt/cli.ini renew -q

[Install]
WantedBy=multi-user.target
