# roles/letsencrypt/templates/ssl-renew.timer
[Unit]
Description=Renew SSL certs every two weeks

[Timer]
OnBootSec=15min
OnUnitActiveSec=2week

[Install]
WantedBy=timers.target
